# .config/autostart/redshift.desktop
[redshift]
Type=Application
Name=Redshift
Exec="$(redshift)"
